# .config/autostart/delete-files.desktop
[Desktop Entry]
Exec=/var/home/user/.local/bin/delete-files
Icon=application-x-shellscript
Name=delete-files
Type=Application
